# main/droid-juicer/droid-juicer.initd
#!/sbin/openrc-run

name="droid-juicer"
description="Android vendor firmware extraction"

depend() {
  after mount
  after tmpfs
  need localmount
}

start_pre() {
  echo "Extracting binary firmware..."
}

start() {
  if [ ! -f /var/lib/droid-juicer/status.json ]; then
    /usr/bin/droid-juicer
    /sbin/reboot
  fi
}

stop() {
  return 0
}

