# modem/q6meowd/q6meowd.initd
#!/sbin/openrc-run
supervisor=supervise-daemon

name="q6meowd"
description="Enable q6voice audio when call is performed with ModemManager"

# Note: q6voice_card/q6voice_device need to be set in /etc/conf.d/q6voiced
command="/usr/bin/q6meowd"
command_args=""
supervise_daemon_args="--user root"

depend() {
	need dbus
}
